# backup/crontabs
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*       *       *       *       *       /backup/backup_db.sh
